# autostart/telegram.desktop
[Desktop Entry]
Exec=telegram-desktop &
Name=telegram
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
